# deploy/supervisor-autodl.conf
; AutoDL Supervisor 配置（单一真相来源 — AUTODL Runbook §10 复制此文件）
; 安装：cp deploy/supervisor-autodl.conf /etc/supervisor/conf.d/aistudio.conf
; 然后：mkdir -p /root/autodl-tmp/logs && supervisorctl reread && supervisorctl update
; ComfyUI 与 backend 共用 backend/.venv（§8.1 pip install -r comfyui/requirements.txt）

[program:comfyui]
command=/root/autodl-tmp/AIStudio/backend/.venv/bin/python /root/autodl-tmp/ComfyUI/main.py --listen 127.0.0.1 --port 8000
directory=/root/autodl-tmp/ComfyUI
autostart=true
autorestart=true
stderr_logfile=/root/autodl-tmp/logs/comfyui.err.log
stdout_logfile=/root/autodl-tmp/logs/comfyui.out.log
environment=CUDA_VISIBLE_DEVICES="0"

[program:aistudio-backend]
command=/root/autodl-tmp/AIStudio/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 7788
directory=/root/autodl-tmp/AIStudio/backend
autostart=true
autorestart=true
stderr_logfile=/root/autodl-tmp/logs/backend.err.log
stdout_logfile=/root/autodl-tmp/logs/backend.out.log
