# deploy/supervisor-autodl.conf
; AutoDL Supervisor 配置（单一真相来源 — AUTODL Runbook §10 复制此文件）
;
; ⚠️ Cursor / 人工改代码硬约束（2026-07-14 事故：误删 nginx+cloudflared → 公网 Cloudflare 1033）
; - 本文件必须始终包含全部 5 个 [program:] 段，缺一不可：
;     comfyui0 · comfyui1 · aistudio-backend · nginx · cloudflared
; - 禁止整文件覆盖 / 只复制片段后 cp 到 /etc/supervisor/conf.d/aistudio.conf
; - 只允许「在现有文件上改需要的行」（如 COMFYUI_NODES、--database-url、环境变量）
; - 改完后必须 diff 确认 nginx + cloudflared 段仍在，再 supervisorctl update
; - 公网验收：curl -s -o /dev/null -w "%{http_code}\n" https://velora.seele0420.cloud/api/health → 200
;
; 安装：cp deploy/supervisor-autodl.conf /etc/supervisor/conf.d/aistudio.conf
; 然后：mkdir -p /root/autodl-tmp/logs && /usr/bin/supervisorctl -c /etc/supervisor/supervisord.conf reread && update
; ComfyUI 与 backend 共用 backend/.venv（§8.1 pip install -r comfyui/requirements.txt）

[program:comfyui0]
command=/root/autodl-tmp/AIStudio/backend/.venv/bin/python /root/autodl-tmp/ComfyUI/main.py --listen 127.0.0.1 --port 8000 --database-url sqlite:////root/autodl-tmp/ComfyUI/user/comfyui_8000.db
directory=/root/autodl-tmp/ComfyUI
autostart=true
autorestart=true
stderr_logfile=/root/autodl-tmp/logs/comfyui0.err.log
stdout_logfile=/root/autodl-tmp/logs/comfyui0.out.log
environment=CUDA_VISIBLE_DEVICES="0"

[program:comfyui1]
command=/root/autodl-tmp/AIStudio/backend/.venv/bin/python /root/autodl-tmp/ComfyUI/main.py --listen 127.0.0.1 --port 8001 --database-url sqlite:////root/autodl-tmp/ComfyUI/user/comfyui_8001.db
directory=/root/autodl-tmp/ComfyUI
autostart=true
autorestart=true
stderr_logfile=/root/autodl-tmp/logs/comfyui1.err.log
stdout_logfile=/root/autodl-tmp/logs/comfyui1.out.log
environment=CUDA_VISIBLE_DEVICES="1"

[program:aistudio-backend]
command=/root/autodl-tmp/AIStudio/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 7788
directory=/root/autodl-tmp/AIStudio/backend
autostart=true
autorestart=true
stderr_logfile=/root/autodl-tmp/logs/backend.err.log
stdout_logfile=/root/autodl-tmp/logs/backend.out.log
environment=DASHSCOPE_API_KEY="sk-f7e2cf5faa1a49fc9f4013e85f004744",APP_ENV="production",AGENT_MOCK_GENERATION="false",COMFYUI_URL="http://127.0.0.1:8000",COMFYUI_NODES="http://127.0.0.1:8000,http://127.0.0.1:8001"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/root/autodl-tmp/logs/nginx.err.log
stdout_logfile=/root/autodl-tmp/logs/nginx.out.log

[program:cloudflared]
command=/usr/local/bin/cloudflared tunnel run --token eyJhIjoiMWY2MDkyMDcwNGYyNGUzNGU4MGE1MjQzMDg3NjMxZjYiLCJ0IjoiZDFkMTJkNmUtN2NlZi00NWNlLTk0ZTktMTg2N2MyZmRlZDc4IiwicyI6Ik9XTmxOR1l5WTJJdE5HWTFOQzAwTXpNd0xUZzVPVGd0TnpBeE9HTXlOekUwTnpBMCJ9
autostart=true
autorestart=true
stderr_logfile=/root/autodl-tmp/logs/cloudflared.err.log
stdout_logfile=/root/autodl-tmp/logs/cloudflared.out.log
